# systemd-timers-service/backup-mysql-hourly.service
[Unit]
Description=backup MYSQL databases hourly
Documentation="https://github.com/Eideen/backupMysql"

[Service]
Type=oneshot
WorkingDirectory=/volum/@backup/skripts/backupMysql/
ExecStart=/volum/@backup/skripts/backupMysql/db_backup_horly.sh
User=backup
Nice=19
