# systemd-timers-service/backup-mysql-hourly.timer
[Unit]
Description=backup mysql databases timer

[Timer]
OnCalendar=*:48

[Install]
WantedBy=timers.target
