# systemd-timers-service/backup-mysql-daily.timer
[Unit]
Description=backup mysql databases daily

[Timer]
OnCalendar=00:51

[Install]
WantedBy=timers.target
